# src/bios-agent-asset.service.in
[Unit]
Description=bios-agent-asset service
After=malamute.service
Requires=malamute.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/bios/bios-db-ro
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/bios-agent-asset

[Install]
WantedBy=multi-user.target
