# src/fty-asset.service.in
[Unit]
Description=fty-asset service
After=malamute.service bios-db-init.service
Requires=malamute.service bios-db-init.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
EnvironmentFile=-@sysconfdir@/default/bios-db-rw
ExecStart=@prefix@/bin/fty-asset @sysconfdir@/@PACKAGE@/fty-asset.cfg

[Install]
WantedBy=bios.target
